# on4kst-irc-bridge.service
[Unit]
Description=ON4KST IRC Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
ExecStart=/usr/bin/python3 /usr/lib/on4kst-irc-bridge/on4kst_irc_bridge.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
